# root/etc/services.d/shorturl/run
#!/usr/bin/with-contenv bash

exec \
  s6-setuidgid abc \
    /app/shorturl \
      -dbpath /data/urls.db
